# airootfs/etc/systemd/systemsddm.service
[Unit]
Description=LightDM
After=neofetch.service

[Service]
ExecStart=/usr/bin/lightdm
Restart=always

[Install]
WantedBy=multi-user.target
