# airootfs/etc/systemd/system/lightdm.service
[Unit]
Description=LightDM
After=network.target

[Service]
ExecStart=/usr/bin/lightdm
Restart=always

[Install]
WantedBy=multi-user.target
