# airootfs/etc/systemd/system/neofetch.service
[Unit]
Description=neofetch
After=network.target

[Service]
ExecStart=/bin/neofetch
Restart=always

[Install]
WantedBy=multi-user.target
