# airootfs/etc/systemd/system/sddm.service
[Unit]
Description=SDDM
After=network.target

[Service]
ExecStart=/usr/bin/sddm
Restart=always

[Install]
WantedBy=multi-user.target
